# php-apache-front-controller/docker/cron.d/echo
*/1 * * * * root /usr/bin/echo "Hello World" >> /var/log/test.log
