# php-apache-front-controller/docker/cron.d/rotate
*/1 * * * * root /usr/sbin/logrotate /etc/logrotate.d/test --state /var/log/logrotate.status
